# pensieve-deploy/systemd/archive-sync.service
[Unit]
Description=Sync Pensieve archive segments to Storage Box
Documentation=https://github.com/andotherstuff/pensieve
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pensieve
Group=pensieve

# Run the sync script
ExecStart=/home/pensieve/pensieve/pensieve-deploy/scripts/sync-archive.sh

# Environment
EnvironmentFile=-/home/pensieve/pensieve/pensieve-deploy/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pensieve-sync

# Don't start multiple syncs
# (timer won't trigger new run if previous is still running)

