# pensieve-deploy/systemd/pensieve.service
[Unit]
Description=Pensieve Nostr Indexer
Documentation=https://github.com/andotherstuff/pensieve
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=pensieve
Group=pensieve
WorkingDirectory=/home/pensieve/pensieve/pensieve-deploy

# Start docker compose
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down

# Restart policy
Restart=always
RestartSec=10

# Give containers time to stop gracefully
TimeoutStopSec=60

# Environment file
EnvironmentFile=-/home/pensieve/pensieve/pensieve-deploy/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pensieve

[Install]
WantedBy=multi-user.target

